# deploy/systemd/f1-beat.service
[Unit]
Description=F1 Web Celery Beat Scheduler
After=network.target redis.service f1-worker.service
Wants=redis.service

[Service]
Type=simple
User=f1web
Group=www-data
WorkingDirectory=/var/www/f1-web/backend
Environment=PATH=/home/f1web/.local/bin:/var/www/f1-web/backend/.venv/bin
EnvironmentFile=/var/www/f1-web/.env
ExecStartPre=/bin/bash -c 'POETRY_PATH=$$(which poetry 2>/dev/null || echo "/home/f1web/.local/bin/poetry"); if [ ! -x "$$POETRY_PATH" ]; then echo "Poetry not found"; exit 1; fi'
ExecStart=/bin/bash -c 'POETRY_PATH=$$(which poetry 2>/dev/null || echo "/home/f1web/.local/bin/poetry"); $$POETRY_PATH run celery -A app.tasks.celery_app beat --loglevel=info'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 重要：确保只有一个Beat实例运行
LockPersonality=yes

[Install]
WantedBy=multi-user.target