# deploy/systemd/f1-worker.service
[Unit]
Description=F1 Web Celery Worker
After=network.target redis.service f1-api.service
Wants=redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/f1-web/backend
Environment=PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/var/www/f1-web/.env
ExecStart=/root/.cache/pypoetry/virtualenvs/f1-web-backend-gSlLT1KT-py3.12/bin/python -m celery -A app.tasks.celery_app worker --loglevel=info --concurrency=2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# 性能优化
LimitNOFILE=65536
MemoryLimit=512M

[Install]
WantedBy=multi-user.target